# x86-64/dotfiles/systemd/etc/systemd/system/powertop-auto-tune.service
[Unit]
Description=PowerTOP auto tune

[Service]
Type=idle
ExecStart=/usr/sbin/powertop --auto-tune

[Install]
WantedBy=multi-user.target
