# x86-64/dotfiles/qde/etc/systemd/user/obex.service
# Doesn't work

[Unit]
Description=obex bluetooth file transfer daemon

[Service]
Type=simple
ExecStart=/usr/lib/bluetooth/obexd -n

[Install]
WantedBy=multi-user.target
